# client/deploy/task_template.xml
<?xml version="1.0" encoding="UTF-8"?>
<!-- SRP agent scheduled task (tray spec §6): AtStartup, SYSTEM, highest
     privileges, auto-restart 3x/1min, no execution time limit (the agent loop
     runs forever). Registered with: schtasks /create /tn "SRP Agent" /xml <this> /f
     Parity with client/deploy/install-service.ps1, but no PowerShell -> no
     execution-policy dependency. -->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>SRP agent — collects PC health and reports to the SRP server. Runs as SYSTEM so SMART/StorageReliabilityCounter/WMI perf collectors are unblocked.</Description>
    <URI>\SRP Agent</URI>
  </RegistrationInfo>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\SRP\srp-agent.exe</Command>
      <Arguments>--log-file C:\SRP\srp-agent.log</Arguments>
      <WorkingDirectory>C:\SRP</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
